# cellbrowser-codex/com.xen.cellbrowser-angel-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.xen.cellbrowser-angel-loop
  Fires the CellBrowser/Codex loop on qi's ANGEL SCHEDULE — the exact times of
  his Beside Pulses (confirmed 2026-05-31), NOT a fixed interval.
  16 daily ticks. Install: cp to ~/Library/LaunchAgents/ then
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.xen.cellbrowser-angel-loop.plist
  Point ProgramArguments at the actual CellBrowser loop driver once built.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xen.cellbrowser-angel-loop</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>$HOME/bin/xen-cellbrowser-pulse</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>12</integer></dict>
    <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>22</integer></dict>
    <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>33</integer></dict>
    <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>44</integer></dict>
    <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>55</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>12</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>22</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>33</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>44</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>55</integer></dict>
    <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>11</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/qi/logs/cellbrowser-angel-loop.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/qi/logs/cellbrowser-angel-loop.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
